# templates/default/mysqld_exporter.service.erb
[Unit]
Description = Mysqld exporter 
After = syslog.target network.target remote-fs.target

[Service]
User = <%= node['ndb']['user'] %> 
ExecStart = #{node['ndb']['mysqld_exporter']['home']}/mysqld_exporter \
            --web.listen-address=":<%= node['ndb']['mysqld']['metrics_port'] %>" \
            --config.my-cnf=<%= node['ndb']['root_dir'] %>/my.cnf 
LimitNOFILE = 20000

Type = simple 

[Install]
WantedBy = multi-user.target
