# templates/default/ndb_mgmd.service.erb
[Unit]
Description = NDB Mgmt Server (ndb_mgmd)
After = syslog.target network.target remote-fs.target <%= node['install']['systemd']['after'] %>
StartLimitIntervalSec=30
StartLimitBurst=20

[Service]
User = <%= node['ndb']['user'] %>
PIDFile = <%= "#{node['ndb']['log_dir']}" %>/ndb_mgmd.pid
EnvironmentFile = <%= node['ndb']['scripts_dir'] %>/mgmd_env_variables
ExecStart = <%= "#{node['ndb']['base_dir']}/scripts/mgm-server-start.sh"  %>
ExecStop = <%= "#{node['ndb']['base_dir']}/scripts/mgm-server-stop.sh"  %>
ExecReload = <%= "#{node['ndb']['base_dir']}/scripts/mgm-server-restart.sh"  %>
Restart=on-failure

Type = forking

[Install]
WantedBy = multi-user.target
