# templates/default/ndbd.erb
 #!/bin/sh
#
# Startup script for ndbd
#
### BEGIN INIT INFO
# Provides:                     ndb_node-*
# Required-Start:                 
# Should-Start:                 
# Required-Stop:                  
# Default-Start:                2 3 4 5
# Default-Stop:                 0 1 6
# Short-Description:            start and stop ndb_node-*
# Description:                  Start/Stop/Restart NDB Data Node: ndb_node-*
### END INIT INFO

# Variables
INIT_PROG=<%= "#{node[:ndb][:base_dir]}/scripts/ndbd-init.sh" %>
START_PROG=<%= "#{node[:ndb][:base_dir]}/scripts/ndbd-start.sh" %>
STOP_PROG=<%= "#{node[:ndb][:base_dir]}/scripts/ndbd-stop.sh" %>
RESTART_PROG=<%= "#{node[:ndb][:base_dir]}/scripts/ndbd-restart.sh" %>


RETVAL=0


init() {
  echo "Executing : $INIT_PROG"
  $INIT_PROG
  return $?
}


start() {
  echo "Executing : $START_PROG"
  $START_PROG
  return $?
}

stop() {
  echo "Executing : $STOP_PROG"
  $STOP_PROG
  return $?
}

restart() {
  echo "Executing : $RESTART_PROG"
  $RESTART_PROG
  return $?
}


# user-supplied parameter to stop/start/restart process.
case "$1" in
  init)
    init
    ;;
  start)
    start
    ;;
  stop)
    stop
    ;;
  restart)
    restart
    ;;
  reload)
    restart
    ;;
  -h|--help)
    echo ""
    echo "usage: <prog> start|stop|restart"
    echo ""
    exit 0
    ;;
  *)
    echo $"Usage: <prog> {start|stop|restart}"
    exit 1
  esac


exit $RETVAL
