# templates/default/ndbmtd.service.erb
[Unit]
Description = NDB Multi-threaded DataNode (ndbmtd)
After = syslog.target network.target remote-fs.target <%= @deps %>
StartLimitIntervalSec=120
StartLimitBurst=10

[Service]
User = <%= node['ndb']['user'] %>
PIDFile = <%= "#{node['ndb']['log_dir']}" %>/ndb_<%= @node_id %>.pid
ExecStart = <%= "#{node['ndb']['base_dir']}/scripts/ndbd-start.sh" %> 
ExecStop = <%= "#{node['ndb']['base_dir']}/scripts/ndbd-stop.sh" %> 
ExecReload = <%= "#{node['ndb']['base_dir']}/scripts/ndbd-restart.sh" %>
TimeoutSec= <%= node['ndb']['ndbd']['systemctl_timeout_sec'] %>
Restart=on-failure

Type = forking

[Install]
WantedBy = multi-user.target
