# templates/default/ndbmtd.service.erb
[Unit]
Description = NDB Multi-threaded DataNode (ndbmtd)
After = syslog.target network.target remote-fs.target <%= @deps %>
StartLimitIntervalSec=120
StartLimitBurst=10

[Service]
User = <%= node['ndb']['user'] %>
Nice = -20
PIDFile = <%= "#{node['ndb']['log_dir']}" %>/<%= node['ndb']['ndbd']['pidfile'] %>
EnvironmentFile = <%= node['ndb']['scripts_dir'] %>/ndbd_env_variables
ExecStart = <%= "#{node['ndb']['base_dir']}/scripts/ndbd-start.sh" %> 
ExecStop = <%= "#{node['ndb']['base_dir']}/scripts/ndbd-stop.sh" %> 
ExecReload = <%= "#{node['ndb']['base_dir']}/scripts/ndbd-restart.sh" %>
TimeoutSec= <%= node['ndb']['ndbd']['systemctl_timeout_sec'] %>
Restart=on-failure
RestartSec=40
LimitMEMLOCK=infinity
LimitRTPRIO = 99
LimitNICE = -20

Type = forking

[Install]
WantedBy = multi-user.target
