# templates/default/rdrs.service.erb
[Unit]
Description = RonDB Rest API Server 
After = syslog.target network.target remote-fs.target <%= @deps %>
StartLimitIntervalSec=80
StartLimitBurst=10

[Service]

<% if node['ndb']['rdrs']['containerize'].eql? "true" -%>
User = root 
<% else %>
User = <%= node['ndb']['user'] %>
<% end %>

<% if node['ndb']['rdrs']['containerize'].eql? "false" -%>
PIDFile = <%= @pid_file %>
<% end %>

ExecStart = <%= node['ndb']['scripts_dir'] %>/rdrs-start.sh
ExecStop = <%= node['ndb']['scripts_dir'] %>/rdrs-stop.sh
ExecReload = <%= node['ndb']['scripts_dir'] %>/rdrs-restart.sh
TimeoutSec = 0
Restart=on-failure

<% if node['ndb']['rdrs']['containerize'].eql? "false" -%>
Type = forking
<% end %>

[Install]
WantedBy = multi-user.target
